# files/etc/uci-defaults/luci-app-onliner
#!/bin/sh
/etc/init.d/onliner enable
/etc/init.d/onliner start
sleep 1
/etc/init.d/rpcd reload
exit 0
